# ops/launchd/com.fx-codex.briefing.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fx-codex.briefing</string>
  <!-- 00/05/10/.../55分にFXシグナルボードを1回実行する。
       run_exclusive.pyで前回処理との重複を拒否し、Discord通知を1通へ集約する。
       RunAtLoadは付けず、インストール直後の臨時通知を発生させない。 -->
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__FX_ROOT__/tools/run_exclusive.py</string>
    <string>--name</string><string>fx-briefing</string>
    <string>--locks-dir</string><string>__FX_ROOT__/logs/locks</string>
    <string>--</string>
    <string>/bin/zsh</string>
    <string>__FX_ROOT__/scripts/fx_briefing_once.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__FX_ROOT__</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Minute</key><integer>0</integer></dict>
    <dict><key>Minute</key><integer>5</integer></dict>
    <dict><key>Minute</key><integer>10</integer></dict>
    <dict><key>Minute</key><integer>15</integer></dict>
    <dict><key>Minute</key><integer>20</integer></dict>
    <dict><key>Minute</key><integer>25</integer></dict>
    <dict><key>Minute</key><integer>30</integer></dict>
    <dict><key>Minute</key><integer>35</integer></dict>
    <dict><key>Minute</key><integer>40</integer></dict>
    <dict><key>Minute</key><integer>45</integer></dict>
    <dict><key>Minute</key><integer>50</integer></dict>
    <dict><key>Minute</key><integer>55</integer></dict>
  </array>
  <key>ProcessType</key><string>Background</string>
  <key>ThrottleInterval</key><integer>30</integer>
  <key>StandardOutPath</key><string>__FX_ROOT__/logs/launchd/briefing.out.log</string>
  <key>StandardErrorPath</key><string>__FX_ROOT__/logs/launchd/briefing.err.log</string>
</dict>
</plist>
