# ops/launchd/com.fx-codex.health.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fx-codex.health</string>
  <!-- 5分毎に時間足別統合ブリーフィングの運用データ鮮度を監視し、
       停止をDiscordへWARNING/CRITICAL通知する。旧fusion経路は判定対象に含めない。 -->
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__FX_ROOT__/tools/run_exclusive.py</string>
    <string>--name</string><string>fx-health</string>
    <string>--locks-dir</string><string>__FX_ROOT__/logs/locks</string>
    <string>--</string>
    <string>__PYTHON__</string>
    <string>__FX_ROOT__/tools/data_freshness_monitor.py</string>
    <string>--root</string><string>__FX_ROOT__</string>
    <string>--config</string><string>ops/freshness_targets_timeframe.json</string>
  </array>
  <key>WorkingDirectory</key><string>__FX_ROOT__</string>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>ProcessType</key><string>Background</string>
  <key>ThrottleInterval</key><integer>30</integer>
  <key>StandardOutPath</key><string>__FX_ROOT__/logs/launchd/health.out.log</string>
  <key>StandardErrorPath</key><string>__FX_ROOT__/logs/launchd/health.err.log</string>
</dict>
</plist>
