# trader/deploy/com.trader.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  trader 常駐スーパーバイザ（LaunchAgent）。
  ログイン時に起動し、120 秒ごとに watchdog.sh を実行して
  docker compose を up 状態に保つ（Mac mini のサーバー化）。
  __TRADER_DIR__ は bootstrap.sh が実パスに置換する。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trader.supervisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__TRADER_DIR__/deploy/watchdog.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__TRADER_DIR__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>StandardOutPath</key>
  <string>__TRADER_DIR__/logs/supervisor.out.log</string>
  <key>StandardErrorPath</key>
  <string>__TRADER_DIR__/logs/supervisor.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
